# qup.service
[Unit]
Description=Q-up Django Application
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/Q-up/backend
# Load environment variables from a file (more secure)
EnvironmentFile=/etc/environment
# Alternatively, specify a dedicated file: 
# EnvironmentFile=/home/ubuntu/Q-up/prod.env 
ExecStart=/home/ubuntu/Q-up/venv/bin/gunicorn \
          --config /home/ubuntu/Q-up/backend/gunicorn_config.py \
          backend.wsgi:application

# Restart service if it fails
Restart=on-failure
RestartSec=5s

# Create required directories
ExecStartPre=/bin/mkdir -p /var/log/gunicorn
ExecStartPre=/bin/chown -R ubuntu:www-data /var/log/gunicorn

[Install]
WantedBy=multi-user.target 